# airootfs/etc/systemd/system/xdg-user-dirs-update.service
[Unit]
Description=Update XDG user dir configuration (LC_ALL=C xdg-user-dirs-update --force)
Documentation=man:xdg-user-dirs-update(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/xdg-user-dirs-update --force

[Install]
WantedBy=default.target
